# config/update-motd.d/10-armbian-header
#!/bin/bash

echo -e "Welcome to \e[0;91mRetroPie \n"

